# systemd/dirmon.timer
[Unit]
Description=Run dirmon repeatedly

[Timer]
OnBootSec=15min
OnUnitActiveSec=6m

[Install]
WantedBy=timers.target
